# crond/ph_sync_max_data.crontab
# Run Pharbers Synchronize Max_Data job
00 00 * * * sh /home/jeorch/work/shell/ph_sync_max_data.sh
